# s6/svc-openclaw/run
#!/usr/bin/with-contenv bash
# Start the OpenClaw gateway as a long-running service
exec s6-setuidgid abc openclaw gateway run
